# conf/foscam_streamer_supervisor.conf
[program:foscam_streamer]
command = docker-compose up
directory = /home/ubuntu/foscam_streamer
