# conf/foscam_streamer_supervisor.conf
[program:foscam_streamer]
command = docker-compose up
directory = /root/foscam_streamer
